# pibems-addon/rootfs/etc/services.d/pibems/run
#!/usr/bin/with-contenv bashio
set -e

cd /app
exec python -m app.main
